# src/sql-d.ui/sqld.ui.service
[Unit]
Description=sql-d

[Service]
Type=simple
ExecStart=/usr/share/SqlD.UI/SqlD.UI

[Install]
WantedBy=multi-user.target
